# systemd/orca_auto-bot@.service
[Unit]
Description=Orca Auto unified Telegram bot for %i
After=network-online.target
Wants=network-online.target
PartOf=orca_auto-runtime@%i.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/orca_auto
Environment=PYTHONUNBUFFERED=1
Environment=ORCA_AUTO_CONFIG=/home/%i/orca_auto/config/orca_auto.yaml
ExecStart=/home/%i/orca_auto/.venv/bin/python -m orca_auto.flow.telegram_bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
